# flowers_back/crontab
*/5 * * * * /usr/local/bin/python /app/charge.py >> /app/cron.log 2>&1
